# docker/crontab
# Remove sandbox databases untouched for over 60 minutes
*/10 * * * * find /var/www/html/database/sandbox -type f -name "*.db" -mmin +60 -delete 2>/dev/null; find /var/www/html/database/sandbox -type f -name "*.lock" -mmin +60 -delete 2>/dev/null
